# squid/supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
nodaemon=true

[program:squid]
command=/usr/libexec/squid/cache_swap.sh && /usr/sbin/squid -f /etc/squid/squid.conf
