# consul/etc/etc_init_docker-nm.conf
description "Docker daemon and dnsmasq"

start on started docker
stop on stopped docker

post-start exec /sbin/restart network-manager
